# modulo3/nfsroot/etc/rc.local
grep -q am33xx_pwm /sys/devices/bone_capemgr.?/slots || echo am33xx_pwm > /sys/devices/bone_capemgr.?/slots
echo bspwm_P9_14_16 > /sys/devices/bone_capemgr.?/slots

